# backend/systemd/basalam-review.service
[Unit]
Description=Basalam Review Sync Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/basalam-review/backend
EnvironmentFile=/root/basalam-review/backend/.env
ExecStart=/usr/bin/python3 -m app.main worker
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=basalam-review

[Install]
WantedBy=multi-user.target
