# backend/systemd/basalam-review.service
[Unit]
Description=Basalam Review Sync Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/behdashtik-basalam-sync/backend
EnvironmentFile=/root/behdashtik-basalam-sync/backend/.env
ExecStart=/root/behdashtik-basalam-sync/backend/.venv/bin/python -m app.main worker
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=basalam-review

[Install]
WantedBy=multi-user.target
